# bgp/volumes/host/h11/supervisord.conf
[supervisord]
nodaemon=true

[program:ip-add]
command=/bin/ip -6 route add default via 2001:aaaa:ffff:1::3 metric 100
[program:ip]
command=/bin/ip -6 route del default via 2001:aaaa:ffff:1::1
